# target/files/etc/profile.d/local-colorlll.sh
# color-lll initialization

# Skip the rest for noninteractive shells.
[ -z "$PS1" ] && return

alias lll='ls -la --color=auto' 2>/dev/null


