# target/files/etc/profile.d/local-colorlll.csh
#! /bin/csh -f
# color-lll initialization

alias lll 'ls -la --color=auto'


