# target/files/etc/rc.d/init.d/livesys-inst-02
#!/bin/bash
#
# livesys-inst-02: Init script for live image
#
# chkconfig: 345 02 97
# description: Init script for live image.

. /etc/init.d/functions

svc_name=livesys-inst-02

if [ "$1" != "start" ] ; then
    exit 0
fi

/sbin/chkconfig ${svc_name} off

mkdir -p /root/centos-livecd-scripts/from-target
/bin/bash -x /root/centos-livecd-scripts/target/${svc_name}.sh    \
    >/root/centos-livecd-scripts/from-target/${svc_name}.log 2>&1
RETVAL=$?

if [ "$RETVAL" -eq 0 ]; then
    action "Running ${svc_name}" /bin/true
else
    action "Running ${svc_name}" /bin/false
fi

exit $RETVAL
